# systemd/antbot-key-rotation.service
[Unit]
Description=AntBot Encryption Key Rotation Service
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/antbot
Environment=PATH=/opt/antbot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/opt/antbot
EnvironmentFile=/opt/antbot/.env
ExecStart=/opt/antbot/venv/bin/python /opt/antbot/scripts/rotate_keys.py --force --verification

[Install]
WantedBy=multi-user.target 